# transmission-nat-pmp/root/etc/s6-overlay/s6-rc.d/svc-mod-transmission-nat-pmp/run
#!/usr/bin/with-contenv bash

GATEWAY=$(grep nameserver /etc/resolv.conf | awk '{print $2}')
TRANSMISSION_USERNAME=$USER TRANSMISSION_PASSWORD=$PASS /transmission-nat-pmp -gateway $GATEWAY